# packages/bsp/istorenext/fix-ifaces-name.service
[Unit]
Description=Fix network interface names before network starts
DefaultDependencies=no
After=systemd-udev-settle.service
Before=network-pre.target network.target networking.service ifupdown-pre.service istorenext-init-network.service

[Service]
Type=oneshot
ExecStart=/usr/lib/istorenext/fix-ifaces-name
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=sysinit.target
